# infra/aec-web.service
[Unit]
Description=Audit Evidence Compiler Web
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/home/youruser/audit-evidence-compiler
EnvironmentFile=/home/youruser/audit-evidence-compiler/.env
ExecStart=/home/youruser/audit-evidence-compiler/.venv/bin/uvicorn web.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
